# backend/deploy/silliconthink-backend.service
[Unit]
Description=Silicon Think Backend
After=network.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/silliconthink-runtime
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=SPRING_CONFIG_ADDITIONAL_LOCATION=file:/etc/silliconthink/
# 数据库名、密码、JWT 等放这里（不要提交到 Git）
EnvironmentFile=-/etc/silliconthink/backend.env
ExecStart=/usr/lib/jvm/java-17-openjdk-amd64/bin/java -Xms256m -Xmx512m -jar /opt/silliconthink-runtime/app.jar
SuccessExitStatus=143
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
