# backend/deploy/silliconthink-backend.service
[Unit]
Description=Silicon Think Backend
After=network-online.target mysql.service
Wants=network-online.target
# Optional: only when blog storage is on a NAS mount (see deploy/nas-storage.md)
# After=network-online.target mysql.service mnt-nas-blog.mount
# RequiresMountsFor=/mnt/nas-blog

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/silliconthink-runtime
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=SPRING_CONFIG_ADDITIONAL_LOCATION=file:/etc/silliconthink/
# DB / JWT / BLOG_STORAGE_ROOT etc. (do not commit real secrets)
EnvironmentFile=-/etc/silliconthink/backend.env
ExecStart=/usr/lib/jvm/java-17-openjdk-amd64/bin/java -Xms256m -Xmx512m -jar /opt/silliconthink-runtime/app.jar
# Logs only to /var/log/silliconthink/backend.log, not journal
StandardOutput=null
StandardError=null
SuccessExitStatus=143
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
